# install/systemd/rnv-sync-reverb.service
[Unit]
Description=RNV Sync Reverb WebSocket server
After=network.target

[Service]
Type=simple
WorkingDirectory=@APP_DIR@
ExecStart=@PHP@ artisan reverb:start --host=127.0.0.1 --port=8081
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
