# install/systemd/rnv-sync-web.service
[Unit]
Description=RNV Sync web panel
After=network.target

[Service]
Type=simple
WorkingDirectory=@APP_DIR@
ExecStart=@PHP@ artisan serve --host=127.0.0.1 --port=8770
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
